# crontab
0 4 * * *	apt update && apt upgrade -y

*/4 * * * * DISPLAY=:0 xdotool mousemove_relative 1 0
